# launchd/dream-heavy.plist.template
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>            <string>com.HIPPO_USER.dream-heavy</string>
    <key>ProgramArguments</key>
    <array>
        <string>/usr/bin/caffeinate</string>
        <string>-i</string>
        <string>HIPPO_REPO/bin/dream-heavy</string>
    </array>
    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key><string>/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin</string>
    </dict>
    <key>StartCalendarInterval</key>
    <dict>
        <key>Hour</key>   <integer>3</integer>
        <key>Minute</key> <integer>0</integer>
    </dict>
    <key>RunAtLoad</key>        <false/>
    <key>StandardOutPath</key>  <string>HIPPO_HOME/.claude/debug/dream-heavy.log</string>
    <key>StandardErrorPath</key><string>HIPPO_HOME/.claude/debug/dream-heavy.err</string>
    <key>WorkingDirectory</key> <string>HIPPO_REPO</string>
</dict>
</plist>
